# ganeti-noded.init
#!/bin/sh
#
# ganeti-noded	Start/Stop Ganeti Cluster node daemon.
#
# chkconfig:	2345 96 04
# description:	Ganeti Cluster node daemon.
#
# processname:	ganeti-noded
# pidfile:	/var/run/ganeti/ganeti-noded.pid
#
### BEGIN INIT INFO
# Provides:          ganeti-noded
# Required-Start:    $syslog $remote_fs
# Required-Stop:     $syslog $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Ganeti Cluster node daemon
# Description:       Ganeti Cluster node daemon
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

# Get service config
[ -f /etc/sysconfig/ganeti ] && . /etc/sysconfig/ganeti

if [ ! -f /var/lib/ganeti/server.pem ]; then
	echo "Missing configuration file /var/lib/ganeti/server.pem"
	exit 1
fi

start() {
	# Check if the service is already running?
	if [ -f /var/lock/subsys/ganeti-noded ]; then
		msg_already_running "Ganeti Cluster node daemon"
		return
	fi

	msg_starting "Ganeti Cluster node daemon"
	daemon --user USER --pidfile /var/run/ganeti/ganeti-noded.pid /usr/sbin/ganeti-noded $NODED_ARGS
	RETVAL=$?
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/ganeti-noded
}

stop() {
	if [ ! -f /var/lock/subsys/ganeti-noded ]; then
		msg_not_running "Ganeti Cluster node daemon"
		return
	fi

	msg_stopping "Ganeti Cluster node daemon"
	killproc --pidfile /var/run/ganeti/ganeti-noded.pid ganeti-noded
	rm -f /var/lock/subsys/ganeti-noded
}

# HUP causes daemon to reopen logs only
reload() {
	if [ ! -f /var/lock/subsys/ganeti-noded ]; then
		msg_not_running "Ganeti Cluster node daemon"
		RETVAL=7
		return
	fi

	msg_reloading "Ganeti Cluster node daemon"
	killproc ganeti-noded -HUP
	RETVAL=$?
}

condrestart() {
	if [ ! -f /var/lock/subsys/ganeti-noded ]; then
		msg_not_running "Ganeti Cluster node daemon"
		RETVAL=$1
		return
	fi

	stop
	start
}

RETVAL=0
# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart)
	stop
	start
	;;
  try-restart)
	condrestart 0
	;;
  reload|force-reload|flush-logs)
	reload
	;;
  status)
	status --pidfile /var/run/ganeti/ganeti-noded.pid ganeti-noded
	exit $?
	;;
  *)
	msg_usage "$0 {start|stop|restart|try-restart|reload|force-reload|flush-logs|status}"
	exit 3
esac

exit $RETVAL
